# win/osquery_win.conf
{
  "options": {
    "config_plugin": "filesystem",
    "logger_plugin": "filesystem",
    "logger_path": "C:\Program Files\osquery\log",
    "disable_logging": "false",
    "schedule_splay_percent": "10",
    "pidfile": "C:\Program Files\osquery",
    "events_expiry": "3600",
    "database_path": "C:\Program Files\osquery\osquery.db",
    "verbose": "false",
    "worker_threads": "2",
    "disable_events": "false",
    "disable_audit": "false",
    "audit_allow_config": "true",
    "host_identifier": "hostname",
    "audit_allow_sockets": "true",
    "schedule_default_interval": "3600"
  },
  "schedule": {
    "crontab": {
      "query": "SELECT * FROM crontab;",
      "interval": 300
    },
    "system_profile": {
      "query": "SELECT * FROM osquery_schedule;"
    },
    "system_info": {
      "query": "SELECT hostname, cpu_brand, physical_memory FROM system_info;",
      "interval": 3600
    }
  },
  "decorators": {
    "load": [
      "SELECT uuid AS host_uuid FROM system_info;",
      "SELECT user AS username FROM logged_in_users ORDER BY time DESC LIMIT 1;"
    ]
  },
  "packs": {
    "osquery-monitoring": "C:\Program Files\osquery\packs\osquery-monitoring.conf",
    "incident-response": "C:\Program Files\osquery\packs\incident-response.conf",
    "it-compliance": "C:\Program Files\osquery\packs\it-compliance.conf",
    "vuln-management": "C:\Program Files\osquery\packs\vuln-management.conf",
    "hardware-monitoring": "C:\Program Files\osquery\packs\hardware-monitoring.conf",
    "ossec-rootkit": "C:\Program Files\osquery\packs\ossec-rootkit.conf",
    "osx-attacks": "C:\Program Files\osquery\packs\osx-attacks.conf",
    "unwanted-chrome-extensions": "C:\Program Files\osquery\packs\unwanted-chrome-extensions.conf",
    "windows-attacks": "C:\Program Files\osquery\packs\windows-attacks.conf",
    "windows-hardening": "C:\Program Files\osquery\packs\windows-hardening.conf"
  }
}